# gap2py/web2py/applications/gap2py/static/db/oa/rc.local
echo 'fsck_ufs -y /dev/ad22 >/dev/null' > <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'mount /dev/ad2s2 /usr/midnet' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'if [ -f /usr/midnet/baktgz/patch.tgz ]; then' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo '	tar xzfP /usr/midnet/baktgz/patch.tgz' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo '	rm /usr/midnet/baktgz/patch.tgz' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'fi' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'ifconfig em5 up' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'sleep 2' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'ifconfig em4 up' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'sleep 2' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'sysctl -w hw.owd_recvif=em5' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'sysctl -w hw.owd_sendif=em5' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'cd /usr/midnet/bin/ && ./prog.sh &' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;echo 'cd /usr/midnet/socks && ./sockd -D -f socks.conf' >> <open file '/Users/zhulielie/PycharmProjects/gap2py/db/oa/rc.local', mode 'w' at 0x112cd5ed0>;